# deployment/launchd/com.betthat.scheduler.odds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label -->
    <key>Label</key>
    <string>com.betthat.scheduler.odds</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/vato/work/Bet-That_(Proof of Concept)/scheduler_odds.py</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/vato/work/Bet-That_(Proof of Concept)</string>

    <!-- Run at load (system startup) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive - restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/vato/work/Bet-That_(Proof of Concept)/logs/scheduler_odds.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vato/work/Bet-That_(Proof of Concept)/logs/scheduler_odds.stderr.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Process Type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Throttle interval (seconds to wait before restarting after crash) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
